# devops/docker-prod/supervisor.conf
[supervisord]
nodaemon=true

[program:runonce]
command=/etc/supervisor/runonce.sh
autorestart=false 
startsecs=0
redirect_stderr=false 
redirect_stderr=false 